# autostart/xmodmap.desktop
[Desktop Entry]
Name=Xmodmap
Exec=/home/greg/dev/dotfiles/scripts/load-xmodmap.sh
Terminal=false
Type=Application
StartupNotify=false
